# systemd/coolify-staging-notify.timer
[Unit]
Description=Poll Coolify for fresh staging deployments

[Timer]
OnBootSec=2min
OnUnitActiveSec=45s
RandomizedDelaySec=10s

[Install]
WantedBy=timers.target
